# systemd/user/openclaw-gateway-watchdog.timer
[Unit]
Description=Run OpenClaw gateway watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=openclaw-gateway-watchdog.service

[Install]
WantedBy=timers.target
